# conf/supervisor/webssh.ini
[program:webssh]
command=python /home/project/webssh/main.py
autorestart=true
autostart=false
stderr_logfile= /home/project/webssh/logs/webssh_stdout.log
stdout_logfile= /home/project/webssh/logs/webssh_stderr.log
numprocs=1
numprocs_start=0
user=s0nnet
